# docker/cron/crontab
* * * * * cd /var/www/html && php artisan schedule:run --no-ansi >> cron.log 2>&1
